# rootfs_dir_busybox/etc/inittab
#etc/inittab
# 系统启动以后运行 /etc/init.d/rcS这个脚本文件。
::sysinit:/etc/init.d/rcS
# 将 console作为控制台终端，也就是 ttymxc0。
console::askfirst:-/bin/sh
# 重启的话运行 /sbin/init。
::restart:/sbin/init
# 按下 ctrl+alt+del组合键的话就运行 /sbin/reboot，看来 ctrl+alt+del组合键用于重启系统。
::ctrlaltdel:/sbin/reboot
# 关机的时候执行 /bin/umount，也就是卸载各个文件系统。
::shutdown:/bin/umount -a -r
# 关机的时候执行 /sbin/swapoff，也就是关闭交换分区。
::shutdown:/sbin/swapoff -a